# crontab
# comment
*/6 * * * * * echo starting ;sleep 1; echo "hola";echo "hierr" >&2 ;sleep 20; echo "bye"
0,30 * * * * * du -hs /home/bercab/Descargas 
*/3  * * * * * smem -P apcron
