# crontab
# comment
*/2 * * * * * echo "each 2 seconds"
0,30 * * * * * echo "second 1 and 30 of each  minute"
